# deploy/systemd/hub-ttyd.service
[Unit]
Description=TerminalClaw - ttyd browser terminal (per-project tmux)
After=network.target docker.service

[Service]
Type=simple
ExecStart=/usr/bin/ttyd -p 7681 -i 127.0.0.1 -b /terminal -a -W -m 6 -t "disableLeaveAlert=true" bash /opt/terminalclaw/term.sh
Restart=always
RestartSec=2
User=root
Environment=TERM=xterm-256color
# Only kill ttyd itself on stop/restart — leave the tmux server (and any running
# claude) alive, so restarts/`tclaw update` don't drop your live sessions.
KillMode=process

[Install]
WantedBy=multi-user.target
